# reports/etc/cron.d/ipgem-reports
0 8 * * * root /etc/ipgem-reports/get-logs && /usr/local/bin/ipgem-report >/dev/null
