# reports/etc/cron.d/ipgem-reports
5 8 * * * root /usr/bin/chronic /bin/sh -c '/etc/ipgem-reports/get-logs && /usr/bin/ipgem-report'
