# reports/etc/cron.d/ipgem-reports
5 8 * * * root /etc/ipgem-reports/get-logs && /usr/bin/ipgem-report >/dev/null
