# gateway/etc/cron.d/ipgem-gateway
0 * * * * root /usr/bin/ipgem-resolver 1 /var/log/ipgem-gateway /var/log/ipgem-gateway.1 >>/var/log/ipgem-resolver
