# gateway/etc/cron.d/ipgem-gateway
0 * * * * root /usr/bin/ipgem-resolver 1 /var/log/ipgem /var/log/ipgem.1 >>/var/log/ipgem-resolver
